# .config/systemd/user/bslocal.service
[Service]
Type=exec
ExecStart=/usr/bin/env bslocal

[Install]
WantedBy=default.target
