# .config/systemd/user/figma-agent.service
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
ExecStart=/usr/bin/env figma-agent

[Install]
WantedBy=default.target
